# script.timer
[Unit]
Description=Timer for script.service

[Timer]
OnActiveSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
