# script.service
[Unit]
Description=Monitoring daemon

[Service]
ExecStart=/path/to/script.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
